# resources/systemd/nemo_email_interlock_status_report.timer
[Unit]
Description=Email interlock status report to given users

[Timer]
OnCalendar=daily

[Install]
WantedBy=multi-user.target
