# resources/systemd/nemo_email_out_of_time_reservation_notification.timer
[Unit]
Description=Email out of time reservation notification every minute in NEMO

[Timer]
OnCalendar=minutely

[Install]
WantedBy=multi-user.target
