# resources/systemd/nemo_manage_tool_qualifications.timer
[Unit]
Description=Manage tool qualifications every day in NEMO

[Timer]
OnCalendar=daily

[Install]
WantedBy=multi-user.target
