# resources/systemd/nemo_cancel_unused_reservations.timer
[Unit]
Description=Cancel unused reservations every minute in NEMO

[Timer]
OnCalendar=minutely

[Install]
WantedBy=multi-user.target
