# resources/systemd/nemo_email_weekend_access.timer
[Unit]
Description=Send weekend access email notification every hour

[Timer]
OnCalendar=hourly

[Install]
WantedBy=multi-user.target
