# resources/systemd/nemo_create_closure_alerts.service
[Unit]
Description=Create closure alerts in NEMO

[Service]
Type=oneshot
ExecStart=/usr/bin/docker exec nemo django-admin create_closure_alerts
