# resources/systemd/nemo_auto_logout_users.timer
[Unit]
Description=Auto logout users every minute in NEMO

[Timer]
OnCalendar=minutely

[Install]
WantedBy=multi-user.target
