# resources/systemd/nemo_check_update_wait_list.timer
[Unit]
Description=Check and update tool wait list every minute in NEMO

[Timer]
OnCalendar=minutely

[Install]
WantedBy=multi-user.target
