# resources/systemd/nemo_email_reservation_reminders.timer
[Unit]
Description=Send reservation reminder emails every 15 minutes

[Timer]
OnCalendar=*:0,15,30,45

[Install]
WantedBy=multi-user.target
